# roles/result/files/result.service
[Unit]
Description=Result

[Service]
EnvironmentFile=/home/result/.env
WorkingDirectory=/home/result
ExecStart= npm start

Restart=always

[Install]
WantedBy=multi-user.target
